# systemd/offline_arch_mirror/mirror.service
[Unit]
Description=Syncing arch_offline mirror every 24h

[Service]
Type=simple
ExecStart=/usr/bin/make_mirror
#User=root
