# systemd/build-timers/build_guided.service
[Unit]
Description=Builds a guided archinstall version of the official ArchLinux ISO

[Service]
Type=simple
ExecStart=/usr/bin/build_guided
