# systemd/build-timers/build_mirror.service
[Unit]
Description=Builds a offline mirror server official ArchLinux ISO

[Service]
Type=simple
ExecStart=/usr/bin/build_mirror
